# unwake.service
[Unit]
Description=Unwake Service

[Service]
Type=simple
ExecStart=/path/to/unwake

[Install]
WantedBy=default.target
